# layers/meta-testdistro/recipes-installer/sysinstall-pkg/sysinstall-pkg/autoinstall.service.in
[Unit]
Description=Automated system installer
Requires=setup-nv-boot-control.service
Wants=time-sync.target
After=time-sync.target setup-nv-boot-control.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@SBINDIR@/tegra-sysinstall
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
