# layers/meta-testdistro/recipes-core/systemd/systemd-conf/dmcrypt-cleanup.service
[Unit]
Description="dm-crypt cleanup"
DefaultDependencies=no
After=cryptsetup.target
ConditionPathExists=!/etc/initrd-release

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/rm -f /run/crypt/passphrase

[Install]
WantedBy=sysinit.target
