# layers/meta-testdistro/recipes-installer/sysinstall-pkg/sysinstall-autoinstaller/autoinstall.service.in
[Unit]
Description=Automated system installer
Wants=time-sync.target
After=time-sync.target

[Service]
Type=oneshot
RuntimeDirectory=autoinstall
RemainAfterExit=yes
ExecStart=@BINDIR@/autoinstall
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
