# layers/meta-testdistro/recipes-testdistro/overlay-setup/overlay-setup/overlay-setup.service.in
[Unit]
Description=Create directories for overlay mounts
DefaultDependencies=no
RequiresMountsFor=/var/extra/overlays

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@SBINDIR@/overlay-setup

[Install]
WantedBy=local-fs.target
